# quantum-currency/systemd/quantum-healing.service
[Unit]
Description=Quantum Currency Auto-Healing Service
After=quantum-gunicorn.service

[Service]
Type=oneshot
User=quantum_user
Group=quantum_user
WorkingDirectory=/opt/quantum-currency
ExecStart=/opt/quantum-currency/scripts/healing_script.sh