# quantum-currency/systemd/quantum-gunicorn.service
[Unit]
Description=Quantum Currency Gunicorn Service
After=network.target

[Service]
User=quantum_user
Group=quantum_user
WorkingDirectory=/opt/quantum-currency
Environment="QC_ENV=Production"
Environment="QC_LAMBDA_BOUNDS=0.8,1.2"
ExecStart=/usr/local/bin/gunicorn \
    --workers $(nproc) \
    --threads 2 \
    --timeout 90 \
    --bind unix:/tmp/gunicorn.sock \
    src.api.main:app \
    --log-level warning \
    --log-file /var/log/quantum-currency/gunicorn.log \
    --access-logfile /var/log/quantum-currency/access.log
ExecReload=/bin/kill -s HUP $MAINPID
KillSignal=SIGTERM
TimeoutStopSec=30
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target